# chef_files/cookbooks/django_deployment/templates/default/site_supervisor.conf.erb
[program:gunicorn]
command = <%= @app_home %>/scripts/go_production_gunicorn.sh                    ; Command to start app
user = ubuntu                                                         ; User to run as
stdout_logfile = <%= @app_home %>/log/gunicorn.log   ; Where to write log messages
redirect_stderr = true
logfile_maxbytes = 100MB
logfile_backups = 5 ; keep only last five logfiles -- avoid overflowing disk with logs

[program:celery]
command = <%= @app_home %>/scripts/go_production_celery.sh                    ; Command to start app
user = ubuntu                                                         ; User to run as
stdout_logfile = <%= @app_home %>/log/celery.log   ; Where to write log messages
redirect_stderr = true
logfile_maxbytes = 100MB
logfile_backups = 5
process_name=celeryd%(process_num)s
numprocs=3

[supervisord]
logfile_maxbytes = 100MB
logfile_backups = 5

[eventlistener:process_end_listener]
command=<%= @env_home %>python <%= @app_home %>/admin/supervisor_pagerduty.py
events=PROCESS_STATE_EXITED
stderr_logfile = <%= @app_home %>/log/event_handler.log   ; Where to write log messages]
